# docker_conf/supervisor.conf
[supervisord]
nodaemon=true
user=root
loglevel=warn

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/src/uwsgi.ini
